# files/etc/rc.local
#!/bin/sh -e
sleep 10
iw reg set CN
exit 0
